# 99-monitor-hotplug.rules
KERNEL=="card1", SUBSYSTEM=="drm",RUN+="/usr/bin/at -f /home/misthios/udevhelper now"
